# root/etc/uci-defaults/99_packagesync
#!/bin/sh

[ -n "$(uci -q get packagesync.config)" ] || sed -i "s|^config packagesync|config packagesync 'config'|" /etc/config/packagesync
